# systemd/bind9-resolvconf.service.d/iserv-dns-over-tls.conf
[Unit]
# Ensure Unbound resolvconf is set up before
After=unbound-resolvconf.service
