# systemd/unbound-resolvconf.service.d/iserv.conf
[Service]
ExecStartPost=/bin/sh -c 'echo nameserver 127.0.2.1 | /sbin/resolvconf -a lo.unbound.static'
ExecStopPost=/sbin/resolvconf -d lo.unbound.static

